# config/supervisord/apps.conf
[program:video_content_classification]
command = gunicorn main:app --timeout=1000 -c /app/src/gunicorn.conf.py
directory=/app/src/
autostart=true
autorestart=true
redirect_stderr=true



