# search/config/supervisord/apps.conf
[program:video_search]
command = gunicorn main:app -c /app/src/gunicorn.conf.py
directory=/app/src/
autostart=true
autorestart=true
redirect_stderr=true



